# branches/backfire_10.03/package/m-route/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 78
/svn/!svn/ver/4889/branches/backfire_10.03/package/m-route/files/etc/hotplug.d
END
