# whiterussian/package/base-files/default/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 72
/svn/!svn/ver/3398/whiterussian/package/base-files/default/etc/hotplug.d
END
